# deploy/supervisor/citadel.conf
[program:setup-extensions]
directory=/opt/citadel
command=python3 /opt/citadel/setup_extensions.py
autostart=true
autorestart=false
startsecs=0
priority=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tailscaled]
command=/usr/local/bin/citadel-tailscaled.sh
autostart=true
autorestart=true
startsecs=2
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tailscale-setup]
command=/usr/local/bin/citadel-tailscale-setup.sh
autostart=true
autorestart=false
startsecs=0
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=/usr/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
startsecs=2
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:hello-world]
directory=/opt/citadel/hello_world
command=/opt/citadel/hello_world/venv/bin/python /opt/citadel/hello_world/app.py
autostart=true
autorestart=true
startsecs=2
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:caddy]
command=/usr/local/bin/citadel-caddy.sh
autostart=true
autorestart=true
startsecs=2
priority=60
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:citadel-final-scan]
directory=/opt/citadel
command=/usr/local/bin/citadel-initial-scan.sh
autostart=true
autorestart=false
startsecs=0
priority=90
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
