# deploy/supervisor/citadel.conf
[program:tailscaled]
command=/usr/local/bin/citadel-tailscaled.sh
autostart=true
autorestart=true
startsecs=2
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:caddy]
directory=/opt/citadel
command=python3 /opt/citadel/setup_extensions.py --runtime
autostart=true
autorestart=true
startsecs=2
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:tailscale-setup]
directory=/opt/citadel
command=python3 /opt/citadel/setup_extensions.py --tailscale
autostart=true
autorestart=false
startsecs=0
priority=40
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=/usr/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
startsecs=2
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:hello-world]
directory=/opt/citadel/hello_world
command=php -S 0.0.0.0:%(ENV_HELLO_WORLD_PORT)s -t /opt/citadel/hello_world
autostart=true
autorestart=true
startsecs=2
priority=50
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:citadel-scan]
directory=/opt/citadel
command=/usr/local/bin/citadel-initial-scan.sh
autostart=true
autorestart=false
startsecs=0
priority=90
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
