# deploy/citadel.container
[Unit]
Description=CITADEL Dashboard with Tailscale

[Container]
ContainerName=citadel
Image=localhost/citadel:latest

# Tailscale needs TUN device and network capabilities
AddCapability=NET_ADMIN
AddCapability=NET_RAW
AddDevice=/dev/net/tun

# Environment
Environment=TS_HOSTNAME=
Environment=TS_AUTHKEY=

# Ports
PublishPort=8443:443

# Full persistence profile:
# keep system + user state in named volumes so package installs,
# config changes and session data survive container recreation.
# Not persistable as normal volumes: /proc, /sys, /dev (kernel-managed).
Volume=citadel-etc:/etc:Z
Volume=citadel-bin:/bin:Z
Volume=citadel-home:/home:Z
Volume=citadel-opt:/opt:Z
Volume=citadel-root:/root:Z
Volume=citadel-sbin:/sbin:Z
Volume=citadel-srv:/srv:Z
Volume=citadel-usr:/usr:Z
Volume=citadel-var:/var:Z

[Install]
WantedBy=default.target
